# docker/root/etc/s6-overlay/s6-rc.d/openssh/run
#!/usr/bin/with-contenv bash

exec 2>&1 \
    /usr/sbin/sshd -e \
      -f /opt/ssh/sshd_config